# elch/files/etc/uci-defaults/00-install-pkgs.sh
#!/bin/sh

opkg install --force-overwrite /pkgs/*
rm -rf  /pkgs
