# elch/files/etc/uci-defaults/00-prepare-ftp-dirs.sh
#!/bin/sh

mkdir -p /mnt/sda1
mkdir -p /mnt/sdb1
